# launchd/io.geoshift.mihomo.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>io.geoshift.mihomo</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/usr/local/lib/geoshift/mihomo-run.sh</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>GEOSHIFT_ENV_FILE</key>
		<string>@@GEOSHIFT_ENV_FILE@@</string>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/var/log/geoshift/mihomo.log</string>
	<key>StandardErrorPath</key>
	<string>/var/log/geoshift/mihomo.err.log</string>
</dict>
</plist>
